# config/systemd/system/tm-sync.timer
[Unit]
Description=TM sync every 10 minutes

[Timer]
OnCalendar=*-*-* *:01/10:00

[Install]
WantedBy=timers.target

