# config/systemd/system/tm-sync.service
[Unit]
Description=TM sync

[Service]
Type=simple
ExecStart=/root/bin/tm-sync
