# systemd/tm-sum.service
[Unit]
Description=TM sum

[Service]
# defaults to Type=simple
ExecStart=/home/ich/bin/tm-sum -r -w /home/ich/public_html/tm/
